# compute-bot.service
[Unit]
Description=$COMPUTE Volume Bot
After=network.target

[Service]
Type=simple
User=compute-bot
Group=compute-bot

# Working directory
WorkingDirectory=/opt/compute-bot

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONDONTWRITEBYTECODE=1"
Environment="PATH=/opt/compute-bot/venv/bin"

# Command to run (uses bot_config.json and .bot_wallet.enc in working directory)
ExecStart=/opt/compute-bot/venv/bin/python bot.py run

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60s
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=compute-bot

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/compute-bot
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
